# crontab
0 0 * * 6 python3 src/main.py day debug >> /logfile
